# elasticsearch/elasticsearch.conf.erb
# ElasticSearch Service
description     "ElasticSearch"
 
start on started sshd
stop on stopped sshd
 
respawn limit 10 5
 
env ES_HOME=/usr/local/elasticsearch
env ES_MIN_MEM=256m
env ES_MAX_MEM=2g

exec su - elasticsearch -c "/usr/local/elasticsearch/bin/elasticsearch -f -Des.config=/etc/elasticsearch/elasticsearch.yml -Des.path.home=/usr/local/elasticsearch -Des.path.logs=/var/log/elasticsearch -Des.path.data=/var/lib/elasticsearch -Des.path.work=/tmp/elasticsearch"

pre-start script
    logger -is "$UPSTART_JOB" "started su - elasticsearch -c \"/usr/local/elasticsearch/bin/elasticsearch -f -Des.config=/etc/elasticsearch/elasticsearch.yml -Des.path.home=/usr/local/elasticsearch -Des.path.logs=/var/log/elasticsearch -Des.path.data=/var/lib/elasticsearch -Des.path.work=/tmp/elasticsearch\" "
end script

pre-stop script
    logger -is "$UPSTART_JOB" "stopped"
end script

